# scripts/wifi-server.service
[Unit]
Description=Sanhum Robot WiFi Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sanhum
ExecStart=/usr/bin/python3 /home/pi/sanhum/src/wifi_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
